# pc_donation/supervisord.conf
[supervisord]
logfile = /home/pc_donation/logs/supervisord.log
user = root
nodaemon=true

[program:gunicorn]
command=/home/pc_donation/boot.sh

[program:nginx]
command=/home/pc_donation/nginx.sh
